# tools/lard/debian/postinst
# create user and group leanux:
/usr/bin/getent group leanux > /dev/null || /usr/sbin/groupadd -r leanux
/usr/bin/getent passwd leanux > /dev/null || /usr/sbin/useradd -r -d /bin/lard -s /sbin/nologin -g leanux leanux

# set permission on lard init script
test -f /etc/init.d/lard && /bin/chmod 700 /etc/init.d/lard

# set permission on lard config file
test -f /etc/lard/lard.conf && /bin/chmod 660 /etc/lard/lard.conf

# set ownership on lard config dir
test -d /etc/lard && /bin/chown -R root:leanux /etc/lard

# some info to the installing user
/bin/echo "start lard on boot: /usr/sbin/chkconfig --add lard"
/bin/echo "start lard now: /usr/sbin/service lard start"

exit 0
