# tools/lard/ubuntu/lard.service
[Unit]
Description=@LARD_DESCR@
ConditionFileIsExecutable=/usr/bin/lard

[Service]
User=@LARD_USER@
Group=@LARD_USER@
Type=forking
ExecStart=/usr/bin/lard -c /etc/lard/lard.conf -f /var/lib/lard/lard.db -d
TimeoutSec=0
StandardOutput=tty
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
